# systemd/telegram-parser-celery.service
[Unit]
Description=NITIN Celery Worker
After=network.target redis.service postgresql.service telegram-parser.service
Wants=redis.service

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/root/nitinleads
EnvironmentFile=/etc/systemd/system/telegram-parser.env
Environment="PATH=/root/nitinleads/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Даем время Redis и PostgreSQL запуститься
ExecStartPre=/bin/sleep 5

ExecStart=/root/nitinleads/venv/bin/celery -A config worker \
          --loglevel=info \
          --logfile=/root/nitinleads/logs/celery-worker.log \
          --concurrency=2

Restart=always
RestartSec=10
TimeoutStartSec=60

StandardOutput=append:/root/nitinleads/logs/celery-worker.log
StandardError=append:/root/nitinleads/logs/celery-worker-error.log

[Install]
WantedBy=multi-user.target
