# systemd/telegram-parser-celery-beat.service
[Unit]
Description=NITIN Celery Beat Scheduler
After=network.target redis.service postgresql.service telegram-parser-celery.service
Wants=redis.service telegram-parser-celery.service

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/root/nitinleads
EnvironmentFile=/etc/systemd/system/telegram-parser.env
Environment="PATH=/root/nitinleads/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Даем время Celery worker запуститься
ExecStartPre=/bin/sleep 10

ExecStart=/root/nitinleads/venv/bin/celery -A config beat \
          --loglevel=info \
          --logfile=/root/nitinleads/logs/celery-beat.log \
          --pidfile=/tmp/celery-beat.pid \
          --schedule=/tmp/celerybeat-schedule

Restart=always
RestartSec=10
TimeoutStartSec=60

StandardOutput=append:/root/nitinleads/logs/celery-beat.log
StandardError=append:/root/nitinleads/logs/celery-beat-error.log

[Install]
WantedBy=multi-user.target
