# systemd/telegram-parser.service
[Unit]
Description=Telegram Parser SaaS Gunicorn daemon
After=network.target postgresql.service redis.service

[Service]
Type=notify
User=root
Group=root
WorkingDirectory=/root/nitinleads
EnvironmentFile=/etc/systemd/system/telegram-parser.env
Environment="PATH=/root/nitinleads/venv/bin"
ExecStart=/root/nitinleads/venv/bin/gunicorn \
          --workers 3 \
          --bind unix:/root/nitinleads/gunicorn.sock \
          --timeout 120 \
          --access-logfile /root/nitinleads/logs/gunicorn-access.log \
          --error-logfile /root/nitinleads/logs/gunicorn-error.log \
          config.wsgi:application

[Install]
WantedBy=multi-user.target
