# systemd/telegram-monitor.service
[Unit]
Description=NITIN Telegram Monitor
After=network.target redis.service postgresql.service telegram-parser.service
Wants=redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/nitinleads
EnvironmentFile=/etc/systemd/system/telegram-parser.env
Environment="PATH=/root/nitinleads/venv/bin:/usr/local/bin:/usr/bin:/bin"

ExecStart=/root/nitinleads/venv/bin/python manage.py start_monitoring

Restart=always
RestartSec=30
TimeoutStartSec=120

StandardOutput=append:/root/nitinleads/logs/telegram-monitor.log
StandardError=append:/root/nitinleads/logs/telegram-monitor-error.log

[Install]
WantedBy=multi-user.target
